# support/upstart/migratio-worker.conf
start on starting migratio
stop on stopping migratio
