# lib/support/upstart/migration-worker.conf
start on starting migration
stop on stopping migration

